# Scripts/com.piekstra.cr-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent template for cr-daemon. Scripts/install.sh substitutes the
  __EXEC_PATH__ and __LOG_DIR__ placeholders with absolute paths before
  installing this into ~/Library/LaunchAgents/.

  KeepAlive is a bare `true`: launchd relaunches the app whenever it is not
  running — crash, wedge-exit, or an OS-initiated SIGTERM (logout, launchctl
  kill) — so the daemon reliably survives restarts and system events. It must
  NOT be exit-code-gated: AppKit routes both a user "Quit" and an OS SIGTERM
  through applicationWillTerminate -> exit(0), so launchd cannot tell them apart
  by exit code. A real user Quit instead unloads this agent (launchctl bootout),
  which is why exit 0 alone must never be treated as "stay down."
  ThrottleInterval prevents a relaunch storm if the app crashes on startup.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.piekstra.cr-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>__EXEC_PATH__</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- Minimum seconds between (re)launches: guards against crash-loop storms. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- GUI/interactive agent (owns an NSStatusItem); avoid background throttling. -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- cr and gh live in Homebrew's bin; launchd starts with a bare PATH. -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
